# scripts/reporting-service/mlnx-ceph-report.service
[Unit]
Description=MLNX CEPH reporting service
#Wants=ceph.target
After=ceph.target

[Service]
ExecStart=/opt/mellanox/ceph/reporting-service/mlnx-ceph-report.sh

[Install]
WantedBy=multi-user.target
